# private_etc/systemd/system/homelab-cleanup.service
[Unit]
Description=Homelab cleanup - prune images, volumes & journal

[Service]
Type=oneshot
ExecStart=/usr/local/bin/homelab-cleanup.sh
